# systemd/user/fileward.service
[Unit]
Description=Fileward local automation daemon
After=default.target
ConditionPathIsDirectory=%h/Downloads

[Service]
Type=simple
ExecStart=%h/.local/bin/fileward %h/Downloads
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
